# how_linux_work/systemd_listening_example/echo.socket
[Unit]
Description=Echo Socket

[Socket]
ListenStream=22222
Accept=yes
