# how_linux_work/systemd_listening_example/echo@.service
[Unit]
Description=Echo Service

[Service]
ExecStart=-/bin/cat
StandardInput=socket

